# cast-manager/systemd/uxplay-receiver.service.example
[Unit]
Description=Cast Manager AirPlay Receiver (UxPlay)
After=network-online.target avahi-daemon.service graphical.target
Wants=network-online.target avahi-daemon.service

[Service]
Type=simple
User=o
Group=o
Environment=DISPLAY=:0
Environment=XAUTHORITY=/run/user/1000/gdm/Xauthority
Environment=XDG_RUNTIME_DIR=/run/user/1000
ExecStart=/usr/bin/uxplay -n "Cast Manager" -nh -vs ximagesink
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
