# cast-manager/cast-manager.service
[Unit]
Description=File Manager - Files, Media, Casting, and Integrated Drive
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=__USER__
Group=__USER__
WorkingDirectory=/home/__USER__/cast_manager_v3
Environment=NODE_ENV=production
Environment=PORT=8004
Environment=PUBLIC_URL=http://REDACTED_SERVER_IP:8004
Environment=CAST_PUBLIC_BASE_URL=http://REDACTED_SERVER_IP:8004
Environment=FILE_MANAGER_LIBRARY=/home/__USER__/file-manager/drive
Environment=FILE_MANAGER_MAX_UPLOAD_MB=4096
Environment=FILE_MANAGER_TEXT_PREVIEW_MB=2
Environment=PATH=/home/__USER__/.nvm/versions/node/v20.20.0/bin:/home/__USER__/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/__USER__/.nvm/versions/node/v20.20.0/bin/node server.js
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cast-manager

[Install]
WantedBy=multi-user.target
